# systemd/weatherdisplay.timer
[Unit]
Description=Refresh the weather display every 10 minutes

[Timer]
OnBootSec=2min
OnUnitActiveSec=10min
AccuracySec=30s
Unit=weatherdisplay.service

[Install]
WantedBy=timers.target
